# systemd/smtp-gate.timer
[Unit]
Description=Refresh SMTP gate rules

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
